# deployment/crm-backend.service
[Unit]
Description=CEM Backend Application
After=network.target
Documentation=https://github.com/dadad132/crm-backend-V2

[Service]
Type=simple
User=crm
Group=crm
WorkingDirectory=/opt/crm-backend
Environment="PATH=/opt/crm-backend/venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/crm-backend/.env

ExecStart=/opt/crm-backend/venv/bin/python start_server.py

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crm-backend

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/crm-backend/data.db /opt/crm-backend/backups /opt/crm-backend/app/uploads /opt/crm-backend/logs

[Install]
WantedBy=multi-user.target
